# pix2light/systemd/pix2light_fadeIn.timer
[Unit]
Description=pix2light fadeIn timer

[Timer]
OnCalendar= *-*-* 08:15:00
Persistent=true

[Install]
WantedBy=timers.target
