# pix2light/systemd/pix2light_screenShot.timer
[Unit]
Description=pix2light nextimage timer

[Timer]
OnCalendar= *-*-* *:0:00
Persistent=true

[Install]
WantedBy=timers.target

