# pix2light/systemd/pix2light_nextImage.timer
[Unit]
Description=pix2light nextimage timer

[Timer]
OnCalendar= *-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target

