# pix2light/systemd/pix2light_shuffle.timer
[Unit]
Description=pix2light shuffle timer

[Timer]
OnCalendar= *-*-* *:0:00
Persistent=true

[Install]
WantedBy=timers.target

