# pix2light/systemd/pix2light_nextImage.service
[Unit]
Description=pix2light nextimage

[Service]
User=artificiel
ExecStart=/usr/bin/sendosc 127.0.0.1 6666 /pix2light/imgLoader/next i 0





